# deploy/systemd/partybox-airplay.service
[Unit]
Description=PartyBox AirPlay Receiver (shairport-sync)
After=network-online.target sound.target pipewire.service wireplumber.service
Wants=network-online.target pipewire.service wireplumber.service

[Service]
Type=simple
User=partybox
Group=partybox
Environment=HOME=/home/partybox
ExecStart=/usr/bin/shairport-sync -c /etc/partybox/shairport-sync.conf --name "PartyBox AirPlay"
Restart=on-failure
RestartSec=2
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
