# deploy/systemd/partybox-airplay.service
[Unit]
Description=PartyBox AirPlay Receiver (shairport-sync)
After=network-online.target sound.target avahi-daemon.service
Wants=network-online.target avahi-daemon.service

[Service]
Type=simple
User=shairport-sync
Group=shairport-sync
ExecStart=/usr/bin/shairport-sync -c /etc/partybox/shairport-sync.conf --name "PartyBox"
Restart=on-failure
RestartSec=4
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
